# IMAGES/debian-bullseye-lxde/supervisord.conf
[supervisord]
nodaemon=true

[program:xvfb]
command=/usr/bin/Xvfb :1 -screen 0 1920x1080x24
autostart=true
autorestart=true
user=dockeruser

[program:x11vnc]
command=/usr/bin/x11vnc -display :1 -nopw -forever -rfbport 5901
autostart=true
autorestart=true
user=dockeruser

[program:startlxde]
command=/usr/bin/startlxde
autostart=true
autorestart=true
user=dockeruser
