# IMAGES/debian-bullseye-gnome/supervisord.conf
[program:xvfb]
command=Xvfb :1 -screen 0 1920x1080x24
autostart=true
autorestart=true
priority=10

[program:x11vnc]
command=x11vnc -display :1 -passwd vncpassword -forever -shared
autostart=true
autorestart=true
priority=20

[program:gnome]
command=/usr/bin/gnome-session
autostart=true
autorestart=true
priority=30
